# wechat-bridge.service
[Unit]
Description=WeChat Claude Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# install-systemd.sh 会自动替换下面两行的路径
WorkingDirectory=PLACEHOLDER_DIR
ExecStart=PLACEHOLDER_DIR/.venv/bin/python PLACEHOLDER_DIR/bridge.py

# 限制 Node.js 堆内存，保护服务器内存
Environment=NODE_OPTIONS=--max-old-space-size=384

# 降低优先级，保护 Caddy / V2Ray 等核心服务
Nice=10

# 崩溃后自动重启，最多等 10 秒
Restart=always
RestartSec=10

# 日志输出到 journald（用 journalctl -u wechat-bridge -f 查看）
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wechat-bridge

[Install]
WantedBy=multi-user.target
